# gen/configs/dnsmasq-instance-mg0
#!/bin/sh
### BEGIN INIT INFO
# Provides:          dnsmasq-instance-mg0
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts dnsmasq instance mg0
# Description:       Starts the dnsmasq DNS server instance mg0
### END INIT INFO

DAEMON=/usr/sbin/dnsmasq
DAEMON_OPTS="-C /etc/dnsmasq-instance-mg0.conf"
NAME=dnsmasq-instance-mg0
DESC="dnsmasq instance mg0"

. /lib/lsb/init-functions

case "$1" in
  start)
    log_daemon_msg "Starting $DESC" "$NAME"
    start-stop-daemon --start --quiet --pidfile /var/run/dnsmasq-instance-mg0.pid --exec $DAEMON -- $DAEMON_OPTS
    log_end_msg $?
    ;;
  stop)
    log_daemon_msg "Stopping $DESC" "$NAME"
    start-stop-daemon --stop --quiet --oknodo --pidfile /var/run/dnsmasq-instance-mg0.pid --exec $DAEMON
    log_end_msg $?
    ;;
  restart|force-reload)
    $0 stop
    sleep 1
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|force-reload}"
    exit 1
esac
